# provision/gollum.cron
#!/bin/sh -e
USER=otto-app
sudo -u $USER /srv/$USER/sync-with-remote-repo.sh
