# provision/gollum.cron
#!/bin/sh -e
*/5 * * * * otto-app /srv/otto-app/sync-with-remote-repo.sh
