# provision/gollum.cron
#!/bin/sh -e
*/5 * * * * otto-app /srv/otto-app/sync-with-remote-repo.sh >> /var/log/monodocs.log 2>&1
*/5 * * * * root /usr/sbin/nginx -s reload >> /var/log/monodocs.log 2>&1
